# supervisor/ptms.conf
# =============================================================================
# Supervisor Configuration for Project Tracking Management System
# =============================================================================
#
# Copy this file to: /etc/supervisor/conf.d/ptms.conf
# Then run:
#   sudo supervisorctl reread
#   sudo supervisorctl update
#   sudo supervisorctl start ptms:*
# =============================================================================

# =============================================================================
# Program Group
# =============================================================================
[group:ptms]
programs=gunicorn,celery-worker,celery-beat

# =============================================================================
# Gunicorn - Django Application Server
# =============================================================================
[program:gunicorn]
command=/var/www/project_tracking/venv/bin/gunicorn config.wsgi:application -c /var/www/project_tracking/backend/gunicorn.conf.py

; Process configuration
directory=/var/www/project_tracking/backend
user=www-data
group=www-data
numprocs=1
autostart=true
autorestart=true
startretries=5
startsecs=10

; Environment variables
environment=DJANGO_SETTINGS_MODULE="config.settings.production",PYTHONPATH="/var/www/project_tracking/backend"

; Logging
stdout_logfile=/var/log/supervisor/ptms_gunicorn.log
stderr_logfile=/var/log/supervisor/ptms_gunicorn_error.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile_backups=10
redirect_stderr=false

; Process management
stopsignal=TERM
stopwaitsecs=30
stopasgroup=true
killasgroup=true

; Priority (lower numbers start first)
priority=10

; =============================================================================
; Celery Worker - Background Task Processor
; =============================================================================
[program:celery-worker]
command=/var/www/project_tracking/venv/bin/celery -A config.celery worker -l info --concurrency=4 -n worker1@%h

; Process configuration
directory=/var/www/project_tracking/backend
user=www-data
group=www-data
numprocs=1
autostart=true
autorestart=true
startretries=5
startsecs=10

; Environment variables
environment=DJANGO_SETTINGS_MODULE="config.settings.production",PYTHONPATH="/var/www/project_tracking/backend"

; Logging
stdout_logfile=/var/log/supervisor/ptms_celery.log
stderr_logfile=/var/log/supervisor/ptms_celery_error.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile_backups=10
redirect_stderr=false

; Process management
stopsignal=TERM
stopwaitsecs=60
stopasgroup=true
killasgroup=true

; Priority
priority=20

; =============================================================================
; Celery Beat - Periodic Task Scheduler
; =============================================================================
[program:celery-beat]
command=/var/www/project_tracking/venv/bin/celery -A config.celery beat -l info --scheduler redbeat.RedBeatScheduler

; Process configuration
directory=/var/www/project_tracking/backend
user=www-data
group=www-data
numprocs=1
autostart=true
autorestart=true
startretries=5
startsecs=10

; Environment variables
environment=DJANGO_SETTINGS_MODULE="config.settings.production",PYTHONPATH="/var/www/project_tracking/backend"

; Logging
stdout_logfile=/var/log/supervisor/ptms_beat.log
stderr_logfile=/var/log/supervisor/ptms_beat_error.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile_backups=10
redirect_stderr=false

; Process management
stopsignal=TERM
stopwaitsecs=30
stopasgroup=true
killasgroup=true

; Priority
priority=30

; =============================================================================
; Optional: Daphne - ASGI/WebSocket Server (if using Channels)
; Uncomment if you need separate WebSocket handling
; =============================================================================
; [program:daphne]
; command=/var/www/project_tracking/venv/bin/daphne -u /tmp/daphne.sock config.asgi:application
; directory=/var/www/project_tracking/backend
; user=www-data
; group=www-data
; autostart=true
; autorestart=true
; stdout_logfile=/var/log/supervisor/ptms_daphne.log
; stderr_logfile=/var/log/supervisor/ptms_daphne_error.log
; environment=DJANGO_SETTINGS_MODULE="config.settings.production"
; priority=15

; =============================================================================
; Event Listeners
; =============================================================================

; Email notifications on process crashes (optional)
; [eventlistener:crashmail]
; command=/usr/local/bin/crashmail -a -m admin@your-domain.com -s '/usr/sbin/sendmail -t -i'
; events=PROCESS_STATE_EXITED,PROCESS_STATE_FATAL
; buffer_size=10
; directory=/tmp
; umask=022
; priority=999
; autostart=true
; autorestart=unexpected
; startsecs=10
; startretries=3
; stopsignal=TERM
; stopwaitsecs=10
; stopasgroup=false
; killasgroup=false

; =============================================================================
; Program for Memory Monitoring (optional)
; Restarts processes if memory usage exceeds limit
; =============================================================================
; [eventlistener:memmon]
; command=/var/www/project_tracking/venv/bin/memmon limit=512MB program=ptms:gunicorn
; events=TICK_60
; buffer_size=10
; directory=/tmp
; umask=022
; priority=999
; autostart=true
; autorestart=unexpected
; startretries=3

; =============================================================================
; Deployment Commands Reference
; =============================================================================
; 
; After making changes to this file:
;   sudo supervisorctl reread
;   sudo supervisorctl update
;
; Start all PTMS services:
;   sudo supervisorctl start ptms:*
;
; Stop all PTMS services:
;   sudo supervisorctl stop ptms:*
;
; Restart all PTMS services:
;   sudo supervisorctl restart ptms:*
;
; Check status:
;   sudo supervisorctl status ptms:*
;
; View logs:
;   sudo tail -f /var/log/supervisor/ptms_gunicorn.log
;   sudo tail -f /var/log/supervisor/ptms_celery.log
;   sudo tail -f /var/log/supervisor/ptms_beat.log
;
; Reload configuration without restart:
;   sudo supervisorctl reload
